# debian/cron.daily
#!/bin/sh
#
# /etc/cron.daily/obnam -- run obnam every day (if enabled)

set -e

ENABLE=no

if [ -e /etc/default/obnam ]
then
    . /etc/default/obnam
fi

if [ "$ENABLE" = yes ]
then
    obnam backup $DIRS
fi
